# backend/text-processor-worker.service
[Unit]
Description=Text Processor Celery Worker
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/text-processor
Environment=FLASK_ENV=production
EnvironmentFile=/home/ubuntu/text-processor/.env

# Start Celery worker
ExecStart=/home/ubuntu/text-processor/venv/bin/celery -A celery_worker.app worker --loglevel=info

# Restart on failure
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
